# docker-rbd-plugin.socket
[Socket]
ListenStream=/run/docker/plugins/rbd.sock

[Install]
WantedBy = sockets.target
